# launchd/com.risingtides.voice-memo-pipeline.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.risingtides.voice-memo-pipeline</string>

    <!--
      TEMPLATE — install.sh regenerates this file with your real $HOME and the
      real Python binary path. Edit install.sh to change defaults; don't edit
      this file by hand (your edits will be overwritten on next install).

      Placeholders substituted by install.sh:
        __HOME__      → your $HOME (e.g. /Users/yourname)
        __REAL_PY__   → the real python3 binary resolved at install time
                        (e.g. /Library/Developer/CommandLineTools/usr/bin/python3)
        __REPO_DIR__  → absolute path to this repo checkout
    -->

    <!--
      Invokes the REAL Python binary (not the /usr/bin/python3 stub that
      re-execs through xcode-select). macOS TCC matches against the real exec
      path, so FDA must be granted to this same path. install.sh resolves and
      substitutes __REAL_PY__ automatically.
    -->
    <!--
      Default scope for launchd-triggered fires: only consider memos recorded
      in the last 7 days. Keeps the watcher cheap and prevents it from
      accidentally grinding the historical backlog when CloudKit shuffles the
      Recordings directory. To backfill older history, run the pipeline
      manually with explicit --backfill --since N.
    -->
    <key>ProgramArguments</key>
    <array>
        <string>__REAL_PY__</string>
        <string>__REPO_DIR__/voice-memo-pipeline.py</string>
        <string>--since</string>
        <string>7</string>
    </array>

    <!-- Fire whenever the Voice Memos CloudKit container changes. -->
    <key>WatchPaths</key>
    <array>
        <string>__HOME__/Library/Group Containers/group.com.apple.VoiceMemos.shared/Recordings</string>
    </array>

    <!-- Coalesce rapid filesystem events: don't refire within 30s of the last run. -->
    <key>ThrottleInterval</key>
    <integer>30</integer>

    <!-- Don't relaunch on its own if it crashes — WatchPaths will trigger it again. -->
    <key>KeepAlive</key>
    <false/>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/voice-memo-pipeline.launchd.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/voice-memo-pipeline.launchd.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>__HOME__/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>LANG</key>
        <string>en_US.UTF-8</string>
    </dict>

    <!-- Limit to one concurrent invocation. -->
    <key>LimitLoadToSessionType</key>
    <string>Aqua</string>
</dict>
</plist>
